# project 4/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn Assignment2Forecast:app -b localhost:5000
directory=/deploy/app
